# debian/brainrotfilter.service
[Unit]
Description=BrainrotFilter YouTube Content Filter
Documentation=https://github.com/IAmDestructoman/brainrotfilter-deb
After=network-online.target squid.service
Wants=network-online.target

[Service]
Type=exec
User=brainrotfilter
Group=brainrotfilter
WorkingDirectory=/opt/brainrotfilter

ExecStart=/opt/brainrotfilter/venv/bin/uvicorn \
    brainrotfilter.analyzer_service:app \
    --host 0.0.0.0 \
    --port 8199 \
    --workers 1 \
    --app-dir /opt/brainrotfilter/src

Restart=on-failure
RestartSec=5

# Environment
# PYTHONPATH must point to src/ (the package root), NOT src/brainrotfilter/
# Setting it to src/brainrotfilter/ caused "from config import config" to
# resolve correctly but broke absolute imports like "import audio_analyzer"
# in parallel_analyzer — those imports resolved to None in thread workers.
Environment=PYTHONPATH=/opt/brainrotfilter/src
Environment=PATH=/opt/brainrotfilter/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=BRAINROT_DB_PATH=/var/lib/brainrotfilter/brainrotfilter.db
Environment=BRAINROT_KEYWORDS_PATH=/etc/brainrotfilter/keywords.json
Environment=BRAINROT_CONFIG_DIR=/etc/brainrotfilter
Environment=BRAINROT_DATA_DIR=/var/lib/brainrotfilter
Environment=BRAINROT_WWW_DIR=/usr/share/brainrotfilter/www
Environment=BRAINROT_LOG_DIR=/var/log/brainrotfilter

# Security hardening
# NOTE: NoNewPrivileges=true prevents sudo from working inside the service.
# Squid config writes are handled by brainrotfilter-squid-apply.service
# (a separate root-level oneshot unit triggered via systemd socket activation)
# rather than calling sudo directly from Python code.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/brainrotfilter /var/log/brainrotfilter /etc/brainrotfilter /tmp
PrivateTmp=false

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
