# debian/postinst
#!/bin/bash
set -e

case "$1" in
    configure)
        # Create brainrotfilter system user if it does not exist
        if ! getent passwd brainrotfilter >/dev/null 2>&1; then
            adduser --system --group --home /var/lib/brainrotfilter \
                --no-create-home --shell /usr/sbin/nologin \
                brainrotfilter || true
        fi

        # Set up Python virtual environment
        if [ ! -d /opt/brainrotfilter/venv ]; then
            echo "Creating Python virtual environment..."
            python3 -m venv /opt/brainrotfilter/venv
        fi

        # Install Python dependencies
        if [ -f /opt/brainrotfilter/requirements.txt ]; then
            echo "Installing Python dependencies..."
            /opt/brainrotfilter/venv/bin/pip install --quiet --upgrade pip
            /opt/brainrotfilter/venv/bin/pip install --quiet \
                -r /opt/brainrotfilter/requirements.txt
        fi

        # Set directory ownership
        chown -R brainrotfilter:brainrotfilter /var/lib/brainrotfilter
        chown -R brainrotfilter:brainrotfilter /var/log/brainrotfilter
        chown -R brainrotfilter:brainrotfilter /etc/brainrotfilter
        chown -R brainrotfilter:brainrotfilter /opt/brainrotfilter

        # Set permissions
        chmod 750 /var/lib/brainrotfilter
        chmod 750 /var/log/brainrotfilter
        chmod 750 /etc/brainrotfilter
        chmod 755 /usr/lib/brainrotfilter/scripts/*.sh

        # Install polkit rule so the brainrotfilter service user can trigger
        # the privileged brainrotfilter-squid-apply.service without sudo.
        # (NoNewPrivileges=true in the main service unit blocks sudo entirely.)
        #
        # polkitd >= 0.106 (Ubuntu 23.04+, Debian bookworm+): JS rules
        # polkitd 0.105   (Ubuntu 22.04, Debian bullseye): pkla local authority
        POLKIT_VER=$(pkaction --version 2>/dev/null | awk '{print $NF}' || echo "0")
        POLKIT_MAJOR=$(echo "$POLKIT_VER" | cut -d. -f1)
        POLKIT_MINOR=$(echo "$POLKIT_VER" | cut -d. -f2)
        if [ "${POLKIT_MAJOR:-0}" -gt 0 ] || [ "${POLKIT_MINOR:-0}" -ge 106 ] 2>/dev/null; then
            # Modern polkit: install JS rules file
            mkdir -p /etc/polkit-1/rules.d
            install -m 0644 /usr/share/brainrotfilter/polkit/50-brainrotfilter.rules \
                /etc/polkit-1/rules.d/50-brainrotfilter.rules
        else
            # Legacy polkit 0.105: install pkla file
            mkdir -p /etc/polkit-1/localauthority/50-local.d
            install -m 0644 /usr/share/brainrotfilter/polkit/50-brainrotfilter.pkla \
                /etc/polkit-1/localauthority/50-local.d/50-brainrotfilter.pkla
        fi

        # Make the privileged helper script executable
        chmod 0755 /usr/lib/brainrotfilter/scripts/apply_squid_config.sh 2>/dev/null || true

        # Initialize Squid SSL certificate database (needed for SSL bump).
        # security_file_certgen must run as the proxy user, which requires
        # the parent /var/lib/squid directory to exist first.
        if [ ! -d /var/lib/squid/ssl_db ]; then
            mkdir -p /var/lib/squid
            chown proxy:proxy /var/lib/squid
            CERTGEN=""
            for p in /usr/lib/squid/security_file_certgen \
                      /usr/lib64/squid/security_file_certgen; do
                [ -x "$p" ] && CERTGEN="$p" && break
            done
            if [ -n "$CERTGEN" ]; then
                "$CERTGEN" -c -s /var/lib/squid/ssl_db -M 16MB 2>/dev/null || true
                chown -R proxy:proxy /var/lib/squid/ssl_db 2>/dev/null || true
            fi
        fi

        # Add brainrotfilter to the proxy group so it can read Squid files
        usermod -aG proxy brainrotfilter 2>/dev/null || true

        # Writable state dir for watchdog / snapshots / backups
        install -d -m 755 /var/lib/brainrotfilter
        chown brainrotfilter:brainrotfilter /var/lib/brainrotfilter || true

        # Squid helper needs to read /etc/brainrotfilter
        chmod 755 /etc/brainrotfilter || true

        # Enable units. The filtering services (icap, watchdog) are gated
        # by a ConditionPathExists on /etc/brainrotfilter/.wizard_complete,
        # so they no-op on boot until the wizard finishes and then start
        # explicitly. Pre-wizard we only run the admin API (brainrotfilter)
        # so the operator can reach the setup wizard.
        systemctl daemon-reload
        systemctl enable brainrotfilter-squid-apply.service || true
        systemctl enable brainrotfilter.service || true
        systemctl enable brainrotfilter-icap.service || true
        systemctl enable brainrotfilter-watchdog.service || true
        systemctl enable brainrotfilter-console.service || true
        systemctl start brainrotfilter.service || true
        # icap/watchdog intentionally not started here — their unit-level
        # ConditionPathExists would skip the start anyway; the wizard's
        # finalize step starts them after writing the flag.

        # SSH disabled by default on appliance — TUI option toggles it
        systemctl disable ssh.service 2>/dev/null || true

        echo ""
        echo "======================================"
        echo "  BrainrotFilter installed!"
        echo "======================================"
        echo ""
        echo "  Admin panel: http://$(hostname -I | awk '{print $1}'):8199"
        echo "  Setup wizard will run on first visit."
        echo ""
        echo "  Service: systemctl status brainrotfilter"
        echo "  Logs:    journalctl -u brainrotfilter -f"
        echo ""
        ;;

    abort-upgrade|abort-remove|abort-deconfigure)
        ;;

    *)
        echo "postinst called with unknown argument: $1" >&2
        exit 1
        ;;
esac

#DEBHELPER#

exit 0
